# Random/bevbot-devops/ansible/roles/immich/templates/immich-ml.container.j2
# ansible/roles/immich/templates/immich-ml.container.j2
[Unit]
Description=Immich Machine Learning

[Container]
Image=ghcr.io/immich-app/immich-machine-learning:${IMMICH_VERSION}
Pod=immich.pod
ContainerName=immich-ml

EnvironmentFile={{ immich_config_dir }}/immich.env
Volume={{ immich_config_dir }}/ml-cache:/cache:z

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
